# deploy/systemd/devforge.service
[Unit]
Description=DevForge (docker compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=%h/devforge
RemainAfterExit=yes
Environment=COMPOSE_PROJECT_NAME=devforge
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d --build
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target


